# dongle/filesystem/wrapper/etc/inittab
#
# inittab
#
::sysinit:/etc/rcS
::restart:/sbin/init
