# dongle/filesystem/tree/etc/inittab
#
# $Id$
#
::sysinit:/etc/rcS
::restart:/sbin/init
